# systemd/http-redirect.service
[Unit]
Description=HTTP to HTTPS Redirect
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/http-redirect.py

# Restart on failure
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=http-redirect

[Install]
WantedBy=multi-user.target
